# debian/etc/systemd/system/oncofix-frontend.service
[Unit]
Description=OncoFix Health Intake Form (Frontend)
After=network.target

[Service]
Type=simple
User=oncofix
Group=oncofix
WorkingDirectory=/opt/oncofix/health-intake-form
ExecStart=/usr/bin/serve -s dist -l 5173
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oncofix-frontend

[Install]
WantedBy=multi-user.target
